# host-server/etc/cron.d/kindle-news-feeds
15,45 * * * * root sh -c "/var/lib/kindle-news-feeds-host/kindle-news-feeds.sh"
